# rootfs/etc/init.d/mega
#!/sbin/openrc-run

extra_commands="healthcheck"

depend() {
    before rclone
}

start_pre() {
    mkdir -p \
        "/chroot/${RC_SVCNAME}/dev" \
        "/chroot/${RC_SVCNAME}/home" \
        "/chroot/${RC_SVCNAME}/tmp"

    mknod -m 0666 "/chroot/${RC_SVCNAME}/dev/null" c 1 3
    mknod -m 0666 "/chroot/${RC_SVCNAME}/dev/urandom" c 1 9

    for CMD in sh; do
        BIN=$(which "$CMD")
        DIR="/chroot/${RC_SVCNAME}${BIN%/*}"
        mkdir -p "$DIR"
        cp "$BIN" "$DIR"
    done

    for BIN in /usr/bin/mega-*; do
        DIR="/chroot/${RC_SVCNAME}${BIN%/*}"
        mkdir -p "$DIR"
        cp "$BIN" "$DIR"
    done

    for LIB in $(ldd /usr/bin/mega-cmd-server | awk '/=>/ {print $3}'); do
        DIR="/chroot/${RC_SVCNAME}${LIB%/*}"
        mkdir -p "$DIR"
        cp "$LIB" "$DIR"
    done

    chown nobody \
        "/chroot/${RC_SVCNAME}/home" \
        "/chroot/${RC_SVCNAME}/tmp"
}

start() {
    ebegin "Starting ${RC_SVCNAME}"

    . /etc/mega/.env

    chnobody /usr/bin/mega-login "$EMAIL" "$PASSWORD"

    ENDPOINT=$(chnobody /usr/bin/mega-webdav --port="100${RC_SVCNAME##*-}" "$DIRECTORY" | awk '{print $NF}')

    echo "$ENDPOINT" > "/chroot/${RC_SVCNAME}/home/endpoint"

    sed -i "s|url = #${RC_SVCNAME}-url|url = ${ENDPOINT}|" /etc/rclone/rclone.conf

    healthcheck && eend 0 || eend 1
}

stop() {
    ebegin "Stopping ${RC_SVCNAME}"

    chnobody /usr/bin/mega-webdav -d --all

    eend $?
}

status() {
    healthcheck && eend 0 || eend 1
}

healthcheck() {
    (wget -q --spider $(cat "/chroot/${RC_SVCNAME}/home/endpoint")) 2> /dev/null
}

chnobody() {
    local CMD="$@"

    su -s /bin/sh -c " \
            HOME=/home \
            unshare -r \
            chroot /chroot/${RC_SVCNAME} \
            ${CMD}
        " nobody
}
