# rootfs/etc/init.d/mega
#!/sbin/openrc-run

extra_commands="healthcheck"

export SANDBOX_NAME="$RC_SVCNAME"
export SANDBOX_NUMBER="${RC_SVCNAME##*-}"
export SANDBOX_PATH="/chroot/${RC_SVCNAME}"

depend() {
    need nftables
    before rclone
}

start_pre() {
    /opt/mega-nfs/sandbox setup_chroot
    /opt/mega-nfs/sandbox setup_network

    for BIN in /bin/sh /usr/bin/mega-*; do
        DIR="${SANDBOX_PATH}${BIN%/*}"
        mkdir -p "$DIR"
        cp "$BIN" "$DIR"
    done

    for LIB in $(ldd /usr/bin/mega-cmd-server | awk '/=>/ {print $3}'); do
        DIR="${SANDBOX_PATH}${LIB%/*}"
        mkdir -p "$DIR"
        cp "$LIB" "$DIR"
    done

    . /etc/mega/.env

    # Accept service traffic
    ip netns exec "$SANDBOX_NAME" nft add rule inet nat prerouting tcp dport "1000${SANDBOX_NUMBER}" dnat ip to "127.0.0.1"
    ip netns exec "$SANDBOX_NAME" nft add rule inet filter input tcp dport "1000${SANDBOX_NUMBER}" accept
    ip netns exec "$SANDBOX_NAME" nft add rule inet filter output tcp sport "1000${SANDBOX_NUMBER}" accept

    # Reject incoming connections that exceed the limits
    ip netns exec "$SANDBOX_NAME" nft add rule inet filter input \
        meter active_streams '{ ip daddr . tcp dport timeout 1s limit rate over 1 bytes/second burst 0 bytes }' \
        meter max_streams "{ ip daddr . tcp dport ct count over ${STREAMS} }" \
        reject with tcp reset

    # Reject outgoing connections that exceed the limits
    ip netns exec "$SANDBOX_NAME" nft add rule inet filter output \
        meter active_transfers '{ ip saddr . tcp sport timeout 1s limit rate over 1 bytes/second burst 0 bytes }' \
        meter max_transfers "{ ip saddr . tcp sport ct count over ${TRANSFERS} }" \
        reject with tcp reset

    # Accept unfiltered incoming connections
    ip netns exec "$SANDBOX_NAME" nft add rule inet filter input accept

    # Accept unfiltered outgoing connections
    ip netns exec "$SANDBOX_NAME" nft add rule inet filter output accept
}

start() {
    ebegin "Starting ${RC_SVCNAME}"

    . /etc/mega/.env

    /opt/mega-nfs/sandbox /usr/bin/mega-login "$EMAIL" "$PASSWORD"

    ENDPOINT=$(
        /opt/mega-nfs/sandbox /usr/bin/mega-webdav --port="1000${SANDBOX_NUMBER}" "$DIRECTORY" \
        | awk '{print $NF}' \
        | sed "s/127.0.0.1/10.0.${SANDBOX_NUMBER}.1/"
    )

    echo "$ENDPOINT" > "${SANDBOX_PATH}/home/endpoint"

    sed -i "s|url = #${RC_SVCNAME}-url|url = ${ENDPOINT}|" /etc/rclone/rclone.conf

    healthcheck && eend 0 || (stop; eend 1)
}

stop() {
    ebegin "Stopping ${RC_SVCNAME}"

    /opt/mega-nfs/sandbox /usr/bin/mega-webdav -d --all

    eend $?
}

status() {
    healthcheck && eend 0 || eend 1
}

healthcheck() {
    exit 0
    (wget -q --spider $(cat "${SANDBOX_PATH}/home/endpoint")) 2> /dev/null
}
