# rootfs/etc/init.d/mega
#!/sbin/openrc-run

. /opt/mega-nfs/jail.sh

extra_commands="healthcheck"

depend() {
    before rclone
}

start_pre() {
    for CMD in sh; do
        BIN=$(which "$CMD")
        DIR="/chroot/${RC_SVCNAME}${BIN%/*}"
        mkdir -p "$DIR"
        cp "$BIN" "$DIR"
    done

    for BIN in /usr/bin/mega-*; do
        DIR="/chroot/${RC_SVCNAME}${BIN%/*}"
        mkdir -p "$DIR"
        cp "$BIN" "$DIR"
    done

    for LIB in $(ldd /usr/bin/mega-cmd-server | awk '/=>/ {print $3}'); do
        DIR="/chroot/${RC_SVCNAME}${LIB%/*}"
        mkdir -p "$DIR"
        cp "$LIB" "$DIR"
    done
}

start() {
    ebegin "Starting ${RC_SVCNAME}"

    . /etc/mega/.env

    jail /usr/bin/mega-login "$EMAIL" "$PASSWORD"

    ENDPOINT=$(jail /usr/bin/mega-webdav --port="100${RC_SVCNAME##*-}" "$DIRECTORY" | awk '{print $NF}')

    echo "$ENDPOINT" > "/chroot/${RC_SVCNAME}/home/endpoint"

    sed -i "s|url = #${RC_SVCNAME}-url|url = ${ENDPOINT}|" /etc/rclone/rclone.conf

    healthcheck && eend 0 || (stop; eend 1)
}

stop() {
    ebegin "Stopping ${RC_SVCNAME}"

    jail /usr/bin/mega-webdav -d --all

    eend $?
}

status() {
    healthcheck && eend 0 || eend 1
}

healthcheck() {
    (wget -q --spider $(cat "/chroot/${RC_SVCNAME}/home/endpoint")) 2> /dev/null
}
