# rootfs/etc/init.d/rclone
#!/sbin/openrc-run

extra_commands="healthcheck"

PID_FILE=/var/run/rclone.pid

depend() {
    before nfs
    need mega-0
    provide portmap rpc.statd
}

start_pre() {
    [ -f "$PID_FILE" ] && return

    mkdir -p \
        /mnt/rclone \
        /var/rclone

    chown nobody \
        /mnt/rclone \
        /var/rclone

    echo "user_allow_other" >> /etc/fuse.conf
}

start() {
    ebegin "Starting ${RC_SVCNAME}"

    . /etc/rclone/.env

    /usr/bin/rclone \
        mount default: /mnt/rclone \
            --config=/etc/rclone/rclone.conf \
            --cache-dir=/var/rclone \
            --check-first \
            --checkers="$STREAMS" \
            --vfs-cache-mode=full \
            --vfs-cache-max-size=$(( $STREAMS * $CHUNK_SIZE * 2 + $TRANSFERS * $CHUNK_SIZE ))M \
            --vfs-read-ahead=$(( $CHUNK_SIZE * 2 ))M \
            --vfs-read-chunk-size="$READ_CHUNK_SIZE"M \
            --vfs-read-chunk-size-limit="$CHUNK_SIZE"M \
            --vfs-refresh \
            --vfs-fast-fingerprint \
            --multi-thread-chunk-size="$CHUNK_SIZE"M \
            --multi-thread-cutoff="$CHUNK_SIZE"M \
            --multi-thread-streams="$STREAMS" \
            --buffer-size="$BUFFER_SIZE"M \
            --transfers="$TRANSFERS" \
            --cutoff-mode=soft \
            --allow-other \
            --ignore-checksum \
            --links \
            --metadata \
            --no-checksum \
            --server-side-across-configs -vv &

    healthcheck && eend 0 || (stop; eend 1)
}

stop() {
    ebegin "Stopping ${RC_SVCNAME}"

    start-stop-daemon --stop --quiet --pidfile "$PID_FILE"

    eend $?
}

status() {
    healthcheck && eend 0 || eend 1
}

healthcheck() {
    mount | grep -q '/mnt/rclone'
}
